# docker/supervisord.conf
; 测试容器
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:init-permissions]
command=bash -c "mkdir -p /www/var/html/storage/logs /www/var/html/storage/app/public /www/var/html/storage/framework/cache /www/var/html/storage/framework/sessions /www/var/html/storage/framework/views /www/var/html/bootstrap/cache && chown -R appuser:appuser /www/var/html/storage /www/var/html/bootstrap/cache && chmod -R 777 /www/var/html/storage /www/var/html/bootstrap/cache"
directory=/www/var/html
autostart=true
autorestart=false
user=root
stdout_logfile=/var/log/supervisor/init-permissions.log
stderr_logfile=/var/log/supervisor/init-permissions.log

[program:build-cache]
command=bash -c "composer run build-cache"
directory=/www/var/html
autostart=true
autorestart=false
user=appuser
stdout_logfile=/var/log/supervisor/build-cache.log
stderr_logfile=/var/log/supervisor/build-cache.log

[program:laravel-serve]
command=php artisan serve --port=80 --host=0.0.0.0
directory=/www/var/html
autostart=true
autorestart=true
user=appuser
stdout_logfile=/var/log/supervisor/laravel-serve.log
stderr_logfile=/var/log/supervisor/laravel-serve.log

[program:queue-worker]
command=php artisan queue:listen -v --sleep=3 --tries=3
directory=/www/var/html
autostart=true
autorestart=true
user=appuser
stdout_logfile=/var/log/supervisor/queue-worker.log
stderr_logfile=/var/log/supervisor/queue-worker.log

[program:schedule-worker]
command=php artisan schedule:work -v
directory=/www/var/html
autostart=true
autorestart=true
user=appuser
stdout_logfile=/var/log/supervisor/schedule-worker.log
stderr_logfile=/var/log/supervisor/schedule-worker.log

[program:mcp-server]
command=php artisan mcp:serve --transport=http 
directory=/www/var/html
autostart=true
autorestart=true
user=appuser
stdout_logfile=/var/log/supervisor/mcp-server.log
stderr_logfile=/var/log/supervisor/mcp-server.log

; [program:log-monitor]
; command=php artisan pail --timeout=0
; directory=/www/var/html
; autostart=true
; autorestart=true
; user=appuser
; stdout_logfile=/var/log/supervisor/log-monitor.log
; stderr_logfile=/var/log/supervisor/log-monitor.log
